# config/systemd/user/nextcloud.automount
[Unit]
Description=Mount WebDAV Service
After=network-online.target
Wants=network-online.target

[Automount]
Where=$HOME/Nextcloud
TimeoutIdleSec=300

[Install]
WantedBy=remote-fs.target
